# scripts/cec-daemon.service.in
[Unit]
Description=CEC Control Daemon
After=multi-user.target
Wants=dbus.service

[Service]
Type=forking
User=@INSTALL_USER@
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/cec-daemon
RemainAfterExit=true
Restart=on-failure

# Set up runtime directory
Environment=XDG_RUNTIME_DIR=/run
RuntimeDirectory=cec-control
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target