# scripts/cec-daemon.service.in
[Unit]
Description=CEC Control Daemon
After=multi-user.target
Wants=dbus.service

[Service]
Type=forking
User=@INSTALL_USER@
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/cec-daemon
PIDFile=/run/cec-control/cec-daemon.pid
Restart=on-failure
RestartSec=5

# Set up runtime directory
Environment=XDG_RUNTIME_DIR=/run
RuntimeDirectory=cec-control
RuntimeDirectoryMode=0755
Environment=CEC_CONTROL_SOCKET=/run/cec-control/socket

[Install]
WantedBy=multi-user.target