# scripts/cec-daemon.service.in
[Unit]
Description=CEC Control Daemon
Documentation=file:///usr/share/doc/cec-control/README.md
After=network.target
Wants=dbus.service

[Service]
Type=simple
User=@INSTALL_USER@
Group=@INSTALL_USER@
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/cec-daemon
Restart=on-failure
RestartSec=5

# Runtime directories and permissions
# Only specify directory name, not full path (systemd adds /run/ prefix)
RuntimeDirectory=cec-control
RuntimeDirectoryMode=0755
StateDirectory=cec-control
StateDirectoryMode=0750
CacheDirectory=cec-control
CacheDirectoryMode=0750
LogsDirectory=cec-control
LogsDirectoryMode=0750

# Set explicit environment variable to indicate we're running as a system service
# This makes detection much more reliable
Environment=CEC_CONTROL_ENVIRONMENT=system_service

# Security hardening
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target